# mysqlhealth.conf
#upstart script for Centos/6 and most other upstart distros.

start on runlevel [345]
stop on starting shutdown
respawn
script
      	exec /root/mysql-health-check/mhc.py
end script

